# ndm_test_package/etc/systemd/system/ndm_test_popov.service
[Service]
ExecStart=/usr/local/bin/ndm_test
StandardOutput=journal
